# scripts/systemd-timers/updateRepositories.service
[Unit]
Description=Update Pacman List

[Service]
Type=simple
ExecStart=/home/netwokz/scripts/update_repositories.sh
